# docker/supervisord.queue.conf
[supervisord]
logfile=/dev/stdout                             ; supervisord log file
logfile_maxbytes=0                              ; maximum size of logfile before rotation
;logfile_backups=10                             ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[inet_http_server]
port=*:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

;programs

[program:queue_worker]
process_name=%(program_name)s_%(process_num)02d
command=python runner.py core.queue_worker run
directory=/usr/src/app/
autostart=true
autorestart=true
stdout_logfile=/usr/src/app/log/queue.log
stdout_logfile_maxbytes=10000
redirect_stderr=true
stopsignal = QUIT

[program:auditor]
process_name=%(program_name)s_%(process_num)02d
command=python runner.py core audit_runner
directory=/usr/src/app/
autostart=true
autorestart=true
stdout_logfile=/usr/src/app/log/queue.log
stdout_logfile_maxbytes=10000
redirect_stderr=true
stopsignal = QUIT
